# piratebox/piratebox/init.d/piratebox
#!/bin/sh
# /etc/init.d/piratebox

### BEGIN INIT INFO
# Provides:          piratebox
# Required-Start:    $all
# Required-Stop:     $network $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      6 1 0
# Short-Description: All services around piratebox
# Description:       Based on /opt/piratebox/conf/piratebox.conf starts:  * DNSMASQ * HOSTAPD * DROOPY * WLAN-config * CHAT
### END INIT INFO

###
#  Default start-stop script for piratebox services on a laptop
#   written by Matthias Strubel  (matthias.strubel@aod-rpg.de) 2011-02-22
#  licenced by gpl, feel free to improve and send me comments


[ -f /etc/default/rcS ] && . /etc/default/rcS

# PATH for /opt piratebox folder
PATH=$PATH:/opt/piratebox/bin


PIRATEBOX=/opt/piratebox
CONF=$PIRATEBOX/conf/piratebox.conf




if [ -f $CONF ] ; then
   . $CONF
else
   echo "Configuration file not found"
   exit 1
fi

#If you can't find this file
if [ -f $PIRATEBOX/conf/init_done ] ; then
   INIT=OK
else
   $PIRATEBOX_FOLDER/bin/install_piratebox.sh init
fi

# Carry out specific functions when asked to by the system
case "$1" in
  start)
	echo "Starting script piratebox "
	$PIRATEBOX/bin/piratebox_modules.sh start
    ;;
  stop)
    echo  "Stopping script piratebox"
	$PIRATEBOX/bin/piratebox_modules.sh stop
    ;;
  restart)    
		$PIRATEBOX/bin/piratebox_modules.sh restart
	;;
  *)
    echo "Usage: /etc/init.d/piratebox {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
