# piratebox/piratebox/rpi/udev/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/piratebox/rpi/bin/starter.sh"
